# Service/etc/systemd/system/quantum-api.service
# /etc/systemd/system/quantum-api.service
[Unit]
Description=QuantumDev Core API (FastAPI/Uvicorn)
After=network-online.target redis-server.service
Wants=network-online.target redis-server.service

[Service]
Type=simple
WorkingDirectory=/root/quantumdev-open

# Env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/root/quantumdev-open
Environment=CHROMA_PERSIST_DIR=/memory/chroma
EnvironmentFile=-/root/quantumdev-open/.env

# Pre-flight: crea dir Chroma e verifica che l'import funzioni
ExecStartPre=/usr/bin/env bash -lc 'mkdir -p /memory/chroma'
ExecStartPre=/root/quantumdev-open/venv/bin/python -c "import backend.quantum_api as m; print(getattr(m, 'BUILD_SIGNATURE','n/a'))"

# Uvicorn (usa il modulo pieno per evitare 'Could not import module "quantum_api"')
ExecStart=/root/quantumdev-open/venv/bin/python -m uvicorn backend.quantum_api:app \
  --host 0.0.0.0 --port 8081 --workers 1

# Restart policy
Restart=always
RestartSec=2
TimeoutStartSec=60
TimeoutStopSec=20
KillMode=mixed
KillSignal=SIGINT

# Limiti/IO
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal

# Utente
User=root

[Install]
WantedBy=multi-user.target
