# Service/etc/systemd/system/quantum-api.service
# /etc/systemd/system/quantum-api.service
[Unit]
Description=QuantumDev Core API (FastAPI/Uvicorn)
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
WorkingDirectory=/root/quantumdev-open
# Variabili ambiente
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/root/quantumdev-open
EnvironmentFile=-/root/quantumdev-open/.env

# Uvicorn
ExecStart=/root/quantumdev-open/venv/bin/uvicorn backend.quantum_api:app \
  --host 0.0.0.0 --port 8081 --workers 1

# Restart policy
Restart=always
RestartSec=3
TimeoutStopSec=15
KillSignal=SIGINT

# (opzionale) limiti/IO
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal

# Utente (se vuoi eseguirlo come root lascialo così)
User=root

[Install]
WantedBy=multi-user.target
