# Service/etc/systemd/system/telegram-bot.service

# /etc/systemd/system/telegram-bot.service
[Unit]
Description=QuantumDev Telegram Bot
After=network-online.target quantum-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/quantumdev-open/scripts
EnvironmentFile=/root/quantumdev-open/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/quantumdev-open/venv/bin/python3 /root/quantumdev-open/scripts/telegram_bot.py
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
